# setup/systemd/jarvis-hub.service
[Unit]
Description=JARVIS event hub (Redis Streams consumer + state.db writer)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
ExecStart=/home/ulrich/Documents/Projects/jarvis/bin/jarvis-hub
Restart=always
RestartSec=2

# Logs moved off /tmp 2026-05-17 per the same rule the voice-agent
# unit applied 2026-05-07: /tmp gets aggressively cleared/rotated,
# losing forensic evidence within hours. New path persists across
# reboots; rotation handled by jarvis-log-rotate.timer.
StandardOutput=append:%h/.local/share/jarvis/logs/hub.log
StandardError=append:%h/.local/share/jarvis/logs/hub.log

Environment=JARVIS_HUB_URL=redis://127.0.0.1:6379

# ────────────────────────────────────────────────────────────────────
# Sandbox hardening (added 2026-05-17 per enterprise plan §P1-SEC-1)
# ────────────────────────────────────────────────────────────────────
# Hub is a Redis Streams consumer + SQLite writer. It needs no audio,
# no video, no sudo, no shell-out — so it can take strict hardening
# the voice-agent can't.
#
# What the hub legitimately needs to write:
#   - %h/.jarvis/hub/state.db (+ -wal, -shm) — canonical state
#   - %h/.local/share/jarvis/logs/hub.log — runtime log
# Everything else can be read-only or denied.

# No privilege escalation (hub never sudos).
NoNewPrivileges=yes

# Per-process /tmp namespace.
PrivateTmp=yes

# Kernel surface lockdown — hub touches none of these.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid

# Block new setuid/setgid binaries; block namespace creation.
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Drop all capabilities. Hub runs as user 'ulrich' so this is
# belt-and-suspenders.
CapabilityBoundingSet=
AmbientCapabilities=

# Restrict address families — hub talks to Redis on TCP (AF_INET +
# AF_INET6) and writes log files via the local fs syscalls.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Filesystem lockdown:
# - ProtectSystem=strict makes /usr /boot /efi /etc all read-only
# - ProtectHome=read-only means ~ is read-only unless explicitly
#   listed in ReadWritePaths.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.jarvis/hub %h/.local/share/jarvis/logs

# Scoped syscall filter — block kernel-module ops, swap, reboot,
# raw IO, debugging primitives. Hub uses none of these legitimately.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged ~@resources ~@reboot ~@module ~@swap ~@debug ~@raw-io
SystemCallErrorNumber=EPERM

[Install]
WantedBy=default.target
