# setup/systemd/jarvis-evolution-nightly.timer
[Unit]
Description=JARVIS self-evolution - periodic AUTO cycle trigger

[Timer]
# Historical unit name, new cadence: check throughout the day. Manual mode skips
# immediately; Auto mode runs the bounded proposal cycle. User-activity and
# deploy-in-flight guards still skip unsafe moments.
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=2min
Unit=jarvis-evolution-nightly.service

[Install]
WantedBy=timers.target
