# setup/systemd/jarvis-backup-local.service
[Unit]
Description=JARVIS — atomic SQLite snapshot of state.db + turn_telemetry.db
Documentation=file:%h/Documents/Projects/jarvis/scripts/jarvis-backup-local.sh
# Don't run during a backup-target's own startup ordering chaos.
After=network.target

[Service]
Type=oneshot
# Pre-create the snapshot dir AND the parent ~/.jarvis dir (mkdir -p
# handles both). Without this, ReadWritePaths below fails the sandbox
# bring-up with status=226/NAMESPACE on first-run / post-deletion
# because systemd refuses to bind-mount a non-existent path. Defense-
# in-depth: install.sh creates the dir too, but a future `rm -rf
# ~/.jarvis/snapshots` shouldn't kill the next backup cycle.
# Mode 700 — snapshots contain conversation + memory content.
ExecStartPre=/bin/mkdir -p %h/.jarvis/snapshots
ExecStartPre=/bin/chmod 700 %h/.jarvis/snapshots
ExecStart=%h/Documents/Projects/jarvis/scripts/jarvis-backup-local.sh

# Snapshot script writes to ~/.jarvis/snapshots/ and reads
# ~/.jarvis/hub/state.db + ~/.local/share/jarvis/turn_telemetry.db.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.jarvis/snapshots
ReadOnlyPaths=%h/.jarvis/hub %h/.local/share/jarvis
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
CapabilityBoundingSet=
AmbientCapabilities=
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallFilter=~@privileged ~@resources ~@reboot ~@module ~@swap ~@debug ~@raw-io
SystemCallErrorNumber=EPERM

[Install]
WantedBy=default.target
