# setup/systemd/jarvis-voice-recycle.service
[Unit]
Description=JARVIS voice-agent graceful recycle (skips if a session is active)
Documentation=file:%h/Documents/Projects/jarvis/bin/jarvis-voice-recycle
# Don't race the agent's own stop/start.
After=jarvis-voice-agent.service

[Service]
Type=oneshot
ExecStartPre=/bin/mkdir -p %h/.local/share/jarvis/logs
ExecStart=/home/ulrich/Documents/Projects/jarvis/bin/jarvis-voice-recycle
# The script can sleep up to ~2x RETRY_SLEEP_S (default 600s) while a
# late-night session is active before it gives up; give headroom before
# systemd considers it hung.
TimeoutStartSec=1200
StandardOutput=append:%h/.local/share/jarvis/logs/voice-recycle.log
StandardError=append:%h/.local/share/jarvis/logs/voice-recycle.log

# Minimal hardening only. This unit must reach the systemd --user manager
# to restart jarvis-voice-agent.service, so the heavier Protect*/PrivateTmp
# sandbox used by other units (which can block systemctl) is deliberately
# omitted.
ProtectKernelTunables=yes
ProtectKernelModules=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# No [Install]: triggered by jarvis-voice-recycle.timer, never at boot.
