# setup/systemd/jarvis-voice-agent.service
[Unit]
Description=JARVIS voice agent (LiveKit worker: STT → LLM → TTS via Groq)
# Needs the SFU to be listening before we try to connect.
After=network-online.target livekit-server.service
Wants=network-online.target
Requires=livekit-server.service

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=120s
WorkingDirectory=%h/Documents/Projects/jarvis/src/voice-agent
# Load LIVEKIT_URL / LIVEKIT_API_KEY / LIVEKIT_API_SECRET from the
# voice-agent's own .env, AND load GROQ_API_KEY from the shared
# cli/.env.local. Two EnvironmentFile directives; later-loaded wins
# on collision, which is fine — they don't overlap in practice.
EnvironmentFile=%h/Documents/Projects/jarvis/src/voice-agent/.env
# Repo-root .env is the centralized LLM key store (GROQ_API_KEY,
# OPENAI_API_KEY, ANTHROPIC_API_KEY, DEEPSEEK_API_KEY — centralized
# 2026-05-15). Without this line every Groq/OpenAI/Anthropic call
# 401s, the LLM cascade falls through to whatever still has creds,
# and Orpheus TTS — which lives behind GROQ_API_KEY — silently swaps
# to EdgeTTS via the FallbackAdapter, so the user hears robotic voice
# instead of Troy. Optional (`-` prefix) so a key-less install starts.
EnvironmentFile=-%h/Documents/Projects/jarvis/.env
EnvironmentFile=%h/Documents/Projects/jarvis/src/cli/.env.local
# Local bridge bearer token (KEY=VALUE format, chmod 600). Sent on
# every voice-agent → bridge HTTP call. Optional file (`-` prefix)
# so the unit doesn't fail to start before P0c is run.
EnvironmentFile=-%h/.jarvis/local-api-token.env
# ── Browser auto-launch DISABLED (2026-05-20) ──────────────────────
# JARVIS was opening browser windows unprompted. The browser subagent
# has THREE paths that launch a browser when the Chrome extension isn't
# connected: the pre_transfer hook (real Chrome), the per-tool-call
# auto-launch in _browser_ext_base.post (real Chrome on a 503), and the
# CDP fallback (bundled Chromium). Extension not connected → all three
# fired, and the CDP backend then threw tool errors in a loop. Live
# failure 2026-05-20 21:37 (handoff "open the driver's project" → real
# Chrome + bundled Chromium + error loop). Disabling all three makes the
# browser subagent bail cleanly with "extension not connected" instead
# of opening anything. REMOVE these once the JARVIS Chrome extension is
# connected and the CDP backend is fixed. (subagents/browser.py + tools/_browser_ext_base.py)
Environment=JARVIS_BROWSER_PRELAUNCH_DISABLE=1
Environment=JARVIS_BROWSER_DISABLE_CDP_FALLBACK=1
Environment=JARVIS_BROWSER_AUTOLAUNCH_DISABLE=1
# Spec 2026-05-24 memory + procedure loop, Track 1 LIVE flip.
# Force-inject a system message into chat_ctx when the user utterance
# matches the save-trigger regex ("remember this" / "save that" /
# "don't forget" / "write this down" / "memorize this") or the recall-
# trigger regex ("do you remember" / "what did I tell you" / "have I
# told you" / "remind me about/of/what"). The supervisor (Sonnet/Haiku)
# is the second gate — it reads the user's full text + the inject and
# decides whether to actually save/recall. Set to "0" to revert to
# shadow mode (matches still log but no inject).
# Spec: docs/superpowers/specs/2026-05-24-jarvis-memory-and-procedure-loop-design.md
Environment=JARVIS_SAVE_TRIGGER_LIVE=1
Environment=JARVIS_RECALL_TRIGGER_LIVE=1
# Spec 2026-05-24 memory + procedure loop, Track 6 (T14 — LAST flip).
# Enable the autonomous reviewer's apply path. After every turn, the
# off-latency llama-3.1-8b-instant reviewer reads the turn and may
# propose kind=skill_create / skill_patch / memory / procedure.
# Without this flag the proposals are logged + discarded (propose-only);
# with =1 the apply path actually writes via skills_authoring (skills)
# or file_memory.add (memory + procedure targets).
# Spec R2 acceptable-junk-rate is <3/day. If junk exceeds that,
# revert to =0 (or unset) and restart.
# Spec: docs/superpowers/specs/2026-05-24-jarvis-memory-and-procedure-loop-design.md
Environment=JARVIS_SKILL_REVIEW_APPLY=1
# Spec 2026-05-24 source-code self-mod, Plane 3 — Step A enable.
# Activates the pattern detector + propose_code_mod voice tool. The
# pattern detector scans turn_telemetry.db every JARVIS_AUTOMOD_PATTERN_INTERVAL_S
# seconds (default 1800) and emits intents to ~/.jarvis/auto-mods/queue.jsonl
# on ≥3-occurrence threshold crossings.
# IMPORTANT: JARVIS_AUTOMOD_SPAWN_LIVE is intentionally NOT set here.
# Without it the spawner is a no-op — intents accumulate in queue.jsonl
# but no subprocess fires. Flip SPAWN_LIVE=1 only after a 7-day queue
# audit confirms the patterns are sane.
# Spec: docs/superpowers/specs/2026-05-24-jarvis-source-code-self-mod-design.md
Environment=JARVIS_AUTOMOD_ENABLED=1
# `start` (not `dev`) — production mode, no file-watch, JSON logs.
# Uses the venv's interpreter directly so we don't need `source
# .venv/bin/activate` in the ExecStart.
ExecStart=%h/Documents/Projects/jarvis/src/voice-agent/.venv/bin/python %h/Documents/Projects/jarvis/src/voice-agent/jarvis_agent.py start
# Auto-recover from crashes / Groq API hiccups / SFU restart.
Restart=always
RestartSec=5
# The agents SDK emits JSON structured logs — append rather than
# overwrite so we can scroll back across restarts during debugging.
# Was /tmp/jarvis-voice-agent.log until 2026-05-07; switched because
# /tmp gets aggressively cleared/rotated (live: lost the 11:57 SFU
# disconnect evidence within 2h, no journal fallback because logs
# go to file not journald). New path persists across reboots.
# Rotation handled by scripts/rotate-jarvis-logs.sh on a daily timer.
StandardOutput=append:%h/.local/share/jarvis/logs/voice-agent.log
StandardError=append:%h/.local/share/jarvis/logs/voice-agent.log

# ────────────────────────────────────────────────────────────────────
# Sandbox hardening (added 2026-05-11)
# ────────────────────────────────────────────────────────────────────
# Defense-in-depth: limit the blast radius if a misbehaving LLM /
# STT mistranscription / prompt injection causes the bash tool to
# run something unintended. These directives are ADDITIVE — none of
# them block the normal voice-agent workload. Specifically NOT set:
#
#   NoNewPrivileges       — would block sudo (kernel rejects setuid
#                            escalation under NoNewPrivileges).
#                            JARVIS needs sudo per CLAUDE.md.
#   ProtectHome           — voice-agent writes to ~/.jarvis and
#                            ~/.local/share/jarvis on every turn.
#   PrivateDevices        — voice-agent needs /dev/snd for audio
#                            and /dev/video* for the webcam tools.
#   ProtectSystem=strict  — would block apt install / sudoers edits
#                            etc. that the user occasionally asks
#                            JARVIS to do. (Skipped to keep "no UX
#                            cost" — re-enable later if desired.)
#   MemoryDenyWriteExecute — PyTorch / ML runtimes use W+X memory.
#   SystemCallFilter      — too risky for incremental hardening.

# Each voice-agent worker gets its own /tmp namespace. Stops a stray
# bash subprocess from reading other processes' temp files (e.g.
# the launcher lock, browser session caches, etc.).
PrivateTmp=yes

# Read-only mounts on stuff JARVIS legitimately never writes to.
# A misbehaving `bash("rm -rf /etc")` becomes EROFS.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Block creation of new setuid/setgid binaries (existing setuid
# binaries like /usr/bin/sudo still work — this only prevents new
# ones from being made, e.g. `chmod u+s /tmp/x` is a no-op).
RestrictSUIDSGID=yes

# Block namespace creation — voice-agent doesn't unshare/chroot.
# RestrictNamespaces intentionally NOT set: it blocks bubblewrap (the per-command bash-tool sandbox) and broke all bash calls 2026-05-20 (JARVIS: "the sandbox is blocking all bash calls"). bwrap IS the per-command sandbox now; this systemd block was redundant AND breaking.

# Block exotic personality flags (e.g. READ_IMPLIES_EXEC) that
# legitimate code never sets.
LockPersonality=yes

# Drop all capabilities except the bare minimum a user-mode process
# already gets. The voice-agent runs as the user 'ulrich' so it had
# no extra caps to begin with — this is a belt-and-suspenders bound
# in case anything in the dep tree tries to acquire one via setcap.
CapabilityBoundingSet=

[Install]
WantedBy=default.target
