# scripts/systemd/jarvis.service
[Unit]
Description=JARVIS Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/Projects/jarvis
ExecStart=/usr/bin/python3 -m src.server.web_server
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1
Environment=JARVIS_HOT_RELOAD=1
Environment=JARVIS_NO_SANDBOX=1
Environment=JARVIS_OWNER=ulrich

[Install]
WantedBy=default.target
