# scripts/systemd/jarvis-autoupdate.timer
[Unit]
Description=JARVIS Auto-Update every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=jarvis-autoupdate.service

[Install]
WantedBy=timers.target
