# scripts/vps/jarvis-deploy-poll.timer
[Unit]
Description=Run jarvis-deploy-poll every 5 minutes

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=20
Persistent=true

[Install]
WantedBy=timers.target
