# os/systemd/jarvis-grpc.service
[Unit]
Description=JARVIS gRPC Core
Documentation=https://github.com/ulrich/jarvis
After=jarvis-brain.service
Requires=jarvis-brain.service

[Service]
Type=simple
User=jarvis
Group=jarvis
ExecStart=/opt/jarvis/bin/jarvis-core
Environment=RUST_LOG=info
Restart=always
RestartSec=3
TimeoutStartSec=15

[Install]
WantedBy=jarvis.target
