# os/systemd/jarvis-memory.service
[Unit]
Description=JARVIS NeuralLattice Memory
Documentation=https://github.com/ulrich/jarvis
After=local-fs.target

[Service]
Type=notify
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis
EnvironmentFile=/opt/jarvis/.env.os
ExecStart=/opt/jarvis/.venv/bin/python -m brain.memory.service
Restart=always
RestartSec=3
WatchdogSec=60
TimeoutStartSec=30

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/jarvis /opt/jarvis
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=jarvis.target
