# os/systemd/jarvis-greeting.service
[Unit]
Description=JARVIS Boot Greeting
Documentation=https://github.com/ulrich/jarvis
After=jarvis-speech.service jarvis-web.service
Requires=jarvis-speech.service

[Service]
Type=oneshot
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis
EnvironmentFile=/opt/jarvis/.env.os
ExecStart=/opt/jarvis/.venv/bin/python -c "\
import asyncio, subprocess;\
async def greet():\
    try:\
        import edge_tts;\
        c = edge_tts.Communicate('Good morning, Ulrich. All systems nominal. JARVIS OS is online.', 'en-US-AndrewMultilingualNeural');\
        await c.save('/tmp/jarvis-greeting.mp3');\
        subprocess.run(['pw-play', '/tmp/jarvis-greeting.mp3'], timeout=15);\
    except Exception as e:\
        print(f'Greeting failed: {e}');\
asyncio.run(greet())"
RemainAfterExit=yes
TimeoutStartSec=30
SupplementaryGroups=audio

[Install]
WantedBy=jarvis.target
