# os/systemd/jarvis-speech.service
[Unit]
Description=JARVIS Speech (STT/TTS)
Documentation=https://github.com/ulrich/jarvis
After=jarvis-brain.service pipewire.service
Requires=jarvis-brain.service
Wants=pipewire.service

[Service]
Type=simple
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis
EnvironmentFile=/opt/jarvis/.env.os
ExecStart=/opt/jarvis/.venv/bin/python -m brain.speech.service
Restart=always
RestartSec=5
TimeoutStartSec=60
SupplementaryGroups=audio

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/jarvis /opt/jarvis
PrivateTmp=true

[Install]
WantedBy=jarvis.target
