# src/os/desktop/systemd/misty-talk.service
[Unit]
Description=misty-talk (voice turn client — listens for wake event, records, STT, think, TTS)
After=misty-core.service pipewire-pulse.socket
Wants=misty-core.service pipewire-pulse.socket

[Service]
Type=simple
WorkingDirectory=%h/jarvis/src/os/desktop
ExecStart=/usr/bin/bun run client/index.ts --auto-listen --url=http://127.0.0.1:8765
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
